# systemd/streamrs.service
[Unit]
Description=streamrs
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/env streamrs
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
